# olhar-angolano.service
[Unit]
Description=Olhar Angolano - Next.js Application
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/olhar-angolano
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/olhar-angolano/data
ReadWritePaths=/var/www/olhar-angolano/logs

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target